# content/etc/profile.d/java.sh
which java &>/dev/null || PATH="$PATH":"$JAVA_HOME"/bin
